# crontab
*/5 * * * * /usr/local/bin/humhub-cron.sh
